# config/templates/systemd/frappe-bench-web.target
[Unit]
After=network.target
Wants={{ aksara_name }}-logica-web.service {{ aksara_name }}-node-socketio.service

[Install]
WantedBy=multi-user.target
